# deployment/systemd/support-automation-digest-hourly.timer
[Unit]
Description=Support Automation — Hourly digest

[Timer]
OnCalendar=hourly
Persistent=true
Unit=support-automation-digest-hourly.service

[Install]
WantedBy=timers.target
